# snf-deploy/files/etc/rc.local
#!/bin/bash

test -e /etc/network/interfaces.synnefo && \
  ifup --interfaces=/etc/network/interfaces.synnefo %VM_PUBLIC_BRIDGE%

test -e /etc/network/interfaces.synnefo && \
  ifup --interfaces=/etc/network/interfaces.synnefo %VM_PRIVATE_BRIDGE%

test -e /disk &&  losetup -j /disk || losetup -f --show /disk

exit 0
