# snf-deploy/files/etc/rc.local
#!/bin/bash

brctl addbr %COMMON_BRIDGE%
ip link set %COMMON_BRIDGE% up

iptables -t mangle -A PREROUTING -i %COMMON_BRIDGE% -p udp -m udp --dport 67 -j NFQUEUE --queue-num 42

if [ %ROUTER_IP% == %NODE_IP% ]; then
  iptables -t nat -A POSTROUTING -o %PUBLIC_IFACE% -s %SUBNET% -j MASQUERADE
  echo 1 > /proc/sys/net/ipv4/ip_forward
  ip addr add %GATEWAY% dev %COMMON_BRIDGE%
  ip route add %SUBNET% dev %COMMON_BRIDGE% src %GATEWAY%
fi

exit 0
